# docker/supervisord.conf
[supervisord]
# Supervisor reste au premier plan pour que Docker garde le conteneur vivant.
nodaemon=true

# Envoie les logs vers stdout/stderr pour que Render puisse les afficher.
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:php-fpm]
# PHP-FPM execute le code Laravel. L'option -F le garde au premier plan.
command=php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
# Nginx recoit les requetes HTTP et transmet les .php a PHP-FPM.
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
